# files/lib/systemd/system/openvpn@pivpn0.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/pivpn-helper.sh start %i
ExecStopPost=/opt/bpi-r2-router-builder/helpers/pivpn-helper.sh stop %i
