# files/lib/systemd/system/openvpn@client.service.d/changes.conf
[Service]
ExecStartPre=/etc/openvpn/vpn-client.sh %i
Restart=
Restart=always
