# files/lib/systemd/system/openvpn@client.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/vpn-client.sh %i
Restart=
Restart=always
