# files/lib/systemd/system/hostapd@ap0.service.d/changes.conf
[Unit]
Requires=sys-subsystem-net-devices-%i.device

[Service]
ExecStopPost=/opt/bpi-r2-router-builder/helpers/ap0-helper.sh
