# files/lib/systemd/system/openvpn@vpn.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/vpn-helper.sh %i
Restart=
Restart=always
