# files/lib/systemd/system/openvpn@vpn.service.d/changes.conf
[Service]
ExecStartPre=/usr/local/bin/helpers/vpn-helper.sh %i
Restart=
Restart=always
