# files/lib/systemd/system/nftables.service.d/changes.conf
[Service]
ExecStart=
ExecStart=/opt/bpi-r2-router-builder/helpers/nftables-script.sh start
ExecReload=
ExecReload=/opt/bpi-r2-router-builder/helpers/nftables-script.sh start
