# files/lib/systemd/system/openvpn@split.service.d/changes.conf
[Service]
ExecStartPre=/etc/openvpn/vpn-split.sh
Restart=
Restart=always
