# files/etc/network/if-up.d/masquerade
#!/bin/bash
if [[ ! -z "${IF_WPA_SSID}" ]]; then
	iptables -t nat -A POSTROUTING -o ${IFACE} -j MASQUERADE
fi
exit 0