# files/etc/network/if-up.d/miniupnpd
#!/bin/bash
if [[ "$IFACE" == "br0" ]]; then
	ip route add 239.0.0.0/8 via 0.0.0.0 dev br0 metric 0
fi
exit 0
