# files/etc/network/if-up.d/pihole
#!/bin/bash

# Load router settings:
[[ -f /etc/default/router-settings ]] && source /etc/default/router-settings

# Are we dealing with the "br0" interface --AND-- option "redirect_dns" is set to "Y"? 
if [[ "$IFACE" == "br0" && "${redirect_dns:-"Y"}" == "N" ]]; then
	# Add rules to redirect IPv4 DNS requests:
	PIHOLE=$(ifconfig br0 | grep 'inet ' | awk '{print $2}')
	nft add rule inet firewall prerouting_lan ip saddr != ${PIHOLE} ip daddr != ${PIHOLE} udp dport 53 dnat to ${PIHOLE}
	nft add rule inet firewall prerouting_lan ip saddr != ${PIHOLE} ip daddr != ${PIHOLE} tcp dport 53 dnat to ${PIHOLE}

	# Add rules to redirect IPv6 DNS requests:
	PIHOLE=$(ifconfig br0 | grep 'inet6 ' | awk '{print $2}')
	nft add rule inet firewall prerouting_lan ip6 saddr != ${PIHOLE} ip6 daddr != ${PIHOLE} udp dport 53 dnat to ${PIHOLE}
	nft add rule inet firewall prerouting_lan ip6 saddr != ${PIHOLE} ip6 daddr != ${PIHOLE} tcp dport 53 dnat to ${PIHOLE}
fi
exit 0
