# files/etc/cron.d/privoxy_blocklist
@reboot root /usr/local/bin/privoxy-blocklist.sh

