# files/etc/cron.d/samba-checker
* * * * * /opt/bpi-r2-router-builder/helpers/usb-helper.sh prep
