# files/etc/cron.d/update_adblock_list
0 0 * * * /usr/local/bin/update_adblock_list

