# files/etc/cron.d/ydns-updater
0,15,30,45 * * * * /opt/ydns-updater/updater.sh > /dev/null
