# files/lib/systemd/system/ssh.service.d/changes.conf
[Service]
ExecStartPre=
ExecStartPre=/opt/bpi-r2-router-builder/helpers/ssh-helper.sh
ExecStartPre=/usr/sbin/sshd -t
