# files/lib/systemd/system/ssh.service.d/changes.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/local/bin/helpers/ssh-helper.sh
ExecStartPre=/usr/sbin/sshd -t
