# files/lib/systemd/system/ddclient.service.d/changes.conf
[Unit]
ConditionFileNotEmpty=/etc/ddclient.conf

[Service]
ExecCondition=bash -c '! cmp --silent /opt/bpi-r2-router-builder/misc/config/ddclient.conf /etc/ddclient.conf'
