# files/lib/systemd/system/bluetooth.service.d/changes.conf
[Service]
ExecStartPre=modprobe stp_chrdev_bt
ExecStart=
ExecStart=/usr/lib/bluetooth/bluetoothd --compat
