# files/lib/systemd/system/wg-quick@.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/wireguard-helper.sh start %i
ExecStopPost=/opt/bpi-r2-router-builder/helpers/wireguard-helper.sh stop %i
