# files/lib/systemd/system/nginx.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/nginx-helper.sh
