# files/lib/systemd/system/hostapd.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/hostapd-helper.sh start hostapd
ExecStopPost=/opt/bpi-r2-router-builder/helpers/hostapd-helper.sh stop hostapd
