# files/etc/network/if-pre-up.d/user-vpn
#!/bin/bash
iptables -A OUTPUT ! -o lo -m owner --uid-owner vpn -j DROP
exit 0
