# files/lib/systemd/system/docker.service.d/changes.conf
[Service]
ExecStartPre=/bin/mkdir -p /var/lib/docker/bin
ExecStartPre=/bin/mkdir -p /var/lib/docker/data
ExecStartPre=/opt/bpi-r2-router-builder/helpers/docker-helper.sh
