# files/lib/systemd/system/docker.service.d/changes.conf
[Service]
ExecStartPre=/usr/local/bin/helpers/docker-helper.sh

