# files/lib/systemd/system/openvpn@freevpn.service.d/changes.conf
[Service]
ExecStartPre=/etc/openvpn/vpn-freevpn.sh
Restart=
Restart=always
