# files/lib/systemd/system/rc-local.service.d/bpiwrt.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helper/router-helper.sh reformat clear
