# files/lib/systemd/system/udev.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/udev-helper.sh
