# files/lib/systemd/system/squid.service.d/changes.conf
[Unit]
Requires=privoxy.service
After=privoxy.service

[Service]
ExecStartPre=
ExecStartPre=/opt/bpi-r2-router-builder/helpers/squid-helper.sh
ExecStartPre=/usr/sbin/squid --foreground -z
