# files/lib/systemd/system/squid.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/squid-helper.sh
