# files/lib/systemd/system/squid.service.d/changes.conf
[Unit]
Requires=privoxy.service
After=privoxy.service
