# files/lib/systemd/system/networking.service.d/changes.conf
[Service]
ExecStartPre=-/sbin/ifup --read-environment eth0
ExecStart=
ExecStart=/sbin/ifup -a --read-environment --exclude=eth0
