# files/lib/systemd/system/networking.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/mt7615-helper.sh
ExecStartPre=-/sbin/ifup eth0
