# files/lib/systemd/system/networking.service.d/changes.conf
[Unit]
Requires=firewall.service
After=firewall.service

[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/networking-helper.sh
ExecStart=
ExecStart=/sbin/ifup -a --read-environment --exclude=eth0
