# files/lib/systemd/system/openvpn@pivpn.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/pivpn-helper.sh start
ExecStopPost=/opt/bpi-r2-router-builder/helpers/pivpn-helper.sh stop
