# files/etc/cron.weekly/privoxy-adblock
@weekly /usr/local/bin/privoxy-blocklist.sh
