# files/lib/systemd/system/minissdpd.service.d/changes.conf
[Service]
ExecCondition=
ExecStart=
ExecStart=/opt/bpi-r2-router-builder/helpers/minissdpd-helper.sh
