# files/lib/systemd/system/smbd.service.d/usb-helper.conf
[Service]
ExecStartPre=/usr/local/bin/helpers/usb-helper.sh prep
