# files/lib/systemd/system/miniupnpd.service.d/changes.conf
[Unit]
BindsTo=sys-subsystem-net-devices-wan.device
After=sys-subsystem-net-devices-wan.device

[Service]
ExecStartPre=
ExecStartPost=
