# files/lib/systemd/system/miniupnpd.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/miniupnpd-helper.sh start
