# files/lib/systemd/system/miniupnpd.service.d/changes.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/local/bin/helpers/miniupnpd-helper.sh start
ExecStopPost=
ExecStopPost=/usr/local/bin/helpers/miniupnpd-helper.sh stop
