# files/lib/systemd/system/miniupnpd.service.d/changes.conf
[Service]
ExecStartPre=
ExecStartPre=/opt/bpi-r2-router-builder/helpers/miniupnpd-helper.sh start
ExecStopPost=
ExecStopPost=/opt/bpi-r2-router-builder/helpers/miniupnpd-helper.sh stop
