# files/lib/systemd/system/clamav-freshclam.service.d/changes.conf
[Service]
ExecStart=
ExecStart=/usr/bin/freshclam -d --foreground=true --on-update-execute=/usr/local/bin/start_clamav.sh
