# files/lib/systemd/system/clamav-freshclam.service.d/changes.conf
[Service]
ExecStart=
ExecStart=/usr/bin/freshclam -d --foreground=true --on-update-execute=/opt/bpi-r2-router-builder/helpers/freshclam-helper.sh

