# files/lib/systemd/system/openvpn@server.service.d/changes.conf
[Service]
ExecStartPre=/etc/openvpn/vpn-server.sh %i
Restart=
Restart=always
