# files/lib/systemd/system/openvpn@server.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/vpn-server.sh %i
Restart=
Restart=always
