# files/lib/systemd/system/privoxy.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/privoxy-helper.sh pre
ExecStartPost=/opt/bpi-r2-router-builder/helpers/privoxy-helper.sh post
