# files/lib/systemd/system/ulogd2.service.d/changes.conf
[Service]
ExecStartPre=/usr/bin/mkfifo /var/log/ulogd/port_trigger.log
