# files/lib/systemd/system/ulogd2.service.d/changes.conf
[Service]
ExecStartPre=/usr/bin/mkfifo /var/run/port_trigger.fifo
