# files/lib/systemd/system/transmission-daemon.service.d/changes.conf
[Service]
EnvironmentFile=-/etc/default/transmission-daemon
User=
User=vpn
Group=
Group=vpn
ExecStartPre=/opt/bpi-r2-router-builder/helpers/transmission-helper.sh start
ExecStart=
ExecStart=/usr/bin/transmission-daemon -f --log-error ${OPTIONS}
ExecStopPost=
ExecStopPost=/opt/bpi-r2-router-builder/helpers/transmission-helper.sh stop
NoNewPrivileges=
NoNewPrivileges=no
Restart=on-failure
RestartSec=10

