# files/lib/systemd/system/transmission-daemon.service.d/changes.conf
[Service]
User=
User=vpn
Group=
Group=vpn
Type=simple
ExecStart=
ExecStart=/usr/bin/transmission-daemon -f --log-error -g /etc/transmission-daemon
Restart=on-failure
RestartSec=10