# files/lib/systemd/system/openvpn@.service.d/changes.conf
[Service]
ExecStartPre=/opt/bpi-r2-router-builder/helpers/openvpn-helper.sh %i
