# files/lib/systemd/system/hostapd@.service.d/changes.conf
[Unit]
BindsTo=sys-subsystem-net-devices-%i.device
After=sys-subsystem-net-devices-%i.device
