# files/lib/systemd/system/hostapd@.service.d/changes.conf
[Unit]
Requires=sys-subsystem-net-devices-%i.device
